# hokidraw-bot.service
[Unit]
Description=Hokidraw 2D Lottery Bot (single instance)
After=network.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=/opt/hokidraw-bot
EnvironmentFile=/opt/hokidraw-bot/.env
ExecStart=/opt/hokidraw-bot/venv/bin/python main.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hokidraw-bot

[Install]
WantedBy=multi-user.target
